# ansible-lab/control-node/playbooks/templates/etc/systemd/system/app.service
[Unit]
Description=app
After=network.target